# scripts/webhook.service
[Unit]
Description=GitHub Webhook Deploy Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/microbubble-agent/scripts/webhook.py
WorkingDirectory=/opt/microbubble-agent
Restart=always
RestartSec=5
EnvironmentFile=-/opt/microbubble-agent/.env.webhook

[Install]
WantedBy=multi-user.target
